# packages/governor_core/deploy/credence-governor-rotate.timer
[Unit]
Description=Weekly raw-capture retention pass (R-D16)

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
